# files/spideroak-init.conf
description "The SpiderOak automatic cloud backup service."

# Start when system starts.
start on runlevel [2345]
stop on runlevel [!2345]

exec su - root -c '/usr/bin/SpiderOak --headless'
